# files/desktop/etc/profile.d/amd.sh
export RADV_PERFTEST=aco,llvm
export AMD_VULKAN_ICD=RADV
